# deploy/vulcan-backup.timer
[Unit]
Description=Daily vulcan-notify database backup

[Timer]
OnCalendar=*-*-* 03:00:00
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
