# deploy/vulcan-backup.service
[Unit]
Description=vulcan-notify database backup

[Service]
Type=oneshot
ExecStart=/opt/vulcan-notify/deploy/backup-db.sh
